# systemd/openclaw-listen.service.in
[Unit]
Description=OpenClaw microphone listener
Documentation=https://github.com/Westlad/listen
Wants=network-online.target
After=network-online.target sound.target

[Service]
Type=simple
User=@USER@
Group=@USER@
SupplementaryGroups=audio
WorkingDirectory=@REPO_ROOT@
EnvironmentFile=@REPO_ROOT@/.env
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=/run/user/@UID@
Environment=PIPEWIRE_RUNTIME_DIR=/run/user/@UID@
Environment=PULSE_SERVER=unix:/run/user/@UID@/pulse/native
ExecStart=@REPO_ROOT@/target/release/listen daemon
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
